# machine/gilgamesh-44/gilgamesh-44.tunnel.service
[Unit]
Description=gilgamesh-44.tunnel
After=network.target
StartLimitIntervalSec=1

[Service]
Type=simple
Restart=always
RestartSec=4
User=gilgamesh
ExecStart=/home/gilgamesh/main.syncthing/monochrome/machine/gilgamesh-44/gilgamesh-44.tunnel

[Install]
WantedBy=multi-user.target
